# docker/config/supervisor/local.conf
; =============================================================================
; Supervisor Local Development Configuration
; Optimized for development workflow with auto-reload
; =============================================================================

[program:octane]
; Use --watch for automatic reloading when files change
; This is the key for seamless local development
command=php -d opcache.enable_cli=1 /var/www/html/artisan octane:frankenphp --host=0.0.0.0 --port=8000 --workers=1 --max-requests=250 --watch
user=laravel
autostart=true
autorestart=true
priority=1
stopsignal=SIGTERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
; Single worker, quick restarts for development
command=php /var/www/html/artisan queue:work --tries=3 --timeout=60 --max-jobs=100 --max-time=1800
process_name=%(program_name)s_%(process_num)02d
numprocs=1
user=laravel
autostart=true
autorestart=true
priority=2
stopsignal=SIGTERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
