# docker/config/supervisor/production.conf
; =============================================================================
; Supervisor Production Configuration
; Maximum performance and reliability
; =============================================================================

[program:octane]
; Production: NO --watch flag, more workers, JIT enabled
command=php -d opcache.enable_cli=1 -d opcache.jit=1255 /var/www/html/artisan octane:frankenphp --host=0.0.0.0 --port=8000 --workers=4 --max-requests=500
user=laravel
autostart=true
autorestart=true
priority=1
stopsignal=SIGTERM
; Longer wait for graceful shutdown in production
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
; Multiple workers for production throughput
command=php /var/www/html/artisan queue:work --tries=3 --timeout=120 --max-jobs=500 --max-time=3600
process_name=%(program_name)s_%(process_num)02d
; Scale workers based on queue load (adjust as needed)
numprocs=2
user=laravel
autostart=true
autorestart=true
priority=2
stopsignal=SIGTERM
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
; Laravel scheduler (runs every minute)
command=/bin/sh -c "while true; do php /var/www/html/artisan schedule:run --verbose --no-interaction; sleep 60; done"
user=laravel
autostart=true
autorestart=true
priority=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
