# docker/config/supervisor/server-fpm.conf
; =============================================================================
; PHP-FPM server program.
; Used by SERVER=fpm, SERVER=nginx, SERVER=caddy.
; Pool config mounted at /usr/local/etc/php-fpm.d/zz-docker.conf via compose.
; =============================================================================

[program:php-fpm]
command=php-fpm -F -R
user=root
autostart=true
autorestart=true
priority=1
stopsignal=SIGQUIT
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
